# profile.d/50_system_paths.sh
export PATH="$HOME/.local/bin:$HOME/.local/share/flutter/bin:$HOME/.cargo/bin/:$PATH"
