# profile.d/99_ble_final.sh

if [ -f  $HOME/src/ble.sh/out/ble.sh ]; then
	source $HOME/src/ble.sh/out/ble.sh
fi

