# profile.d/99_ble.sh

if [ "$(type -t blesh-share)" = "file" ]; then
    source  "$(blesh-share)/ble.sh"
fi
